# scripts/systemd/flexconnectd.service
[Unit]
Description=FlexConnect AnyConnect daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/sbin/flexconnectd --socket /var/run/flexconnect.sock --state /var/lib/flexconnect/state.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
